# deployment/systemd/bo-vpn-runner.service.example
[Unit]
Description=BO/VPN Runner Daemon
Wants=network-online.target
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/home/timur/bo-vpn-agent
EnvironmentFile=/etc/bo-vpn-runner.env
ExecStart=/home/timur/bo-vpn-agent/.venv/bin/bo-vpn-runner-daemon --host 0.0.0.0 --port 8091
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
